# files/system/silverblue/usr/lib/systemd/user/pixi-global-update.service
[Unit]
Description=Update pixi global packages

[Service]
Type=oneshot
ExecStart=/usr/bin/pixi global update
